# btrfs/scripts/install_stage_1.service
[Unit]
Description=Stage 1 of BTRFS Install
After=network.target

[Service]
ExecStart=/bin/bash /opt/fedora-main/btrfs-fed-37/scripts/btrfs-st-1.sh

[Install]
WantedBy=multi-user.target